# crontabs/root
*/1 * * * * /notify.py @all "it's time for our meeting"
